# hypr-smartd/system/hypr-smartd.service
[Unit]
Description=hypr-smartd layout daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/hypr-smartd --config %h/.config/hypr-smartd/config.yaml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
