# system/hyprpal.service
[Unit]
Description=hyprpal layout daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/hyprpal --config %h/.config/hyprpal/config.yaml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
